# cake-att-reset.service
[Unit]
Description=Nightly CAKE Reset for ATT
After=network-online.target

[Service]
Type=oneshot
User=kevin
WorkingDirectory=/home/kevin/adaptive_cake_att
ExecStart=/usr/bin/python3 /home/kevin/adaptive_cake_att/adaptive_cake_att.py --reset
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cake-att-reset

[Install]
WantedBy=multi-user.target
