# cake-spectrum.service
[Unit]
Description=Adaptive CAKE Auto-Tuning for Spectrum
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=kevin
WorkingDirectory=/home/kevin/adaptive_cake_spectrum
ExecStart=/usr/bin/python3 /home/kevin/adaptive_cake_spectrum/adaptive_cake_spectrum.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cake-spectrum

[Install]
WantedBy=multi-user.target
