# deploy/systemd/cake-autorate-spectrum-state-bridge.service
[Unit]
Description=Bridge Spectrum cake-autorate summary into wanctl state
After=cake-autorate-spectrum.service
Wants=cake-autorate-spectrum.service

[Service]
Type=simple
Environment=WANCTL_EXTERNAL_WAN_NAME=spectrum
Environment=WANCTL_EXTERNAL_DL_IF=spec-router
Environment=WANCTL_EXTERNAL_UL_IF=spec-modem
Environment=CAKE_AUTORATE_BRIDGE_LOG=/var/log/cake-autorate/cake-autorate.spectrum.log
Environment=WANCTL_EXTERNAL_STATE_PATH=/var/lib/wanctl/spectrum_state.json
Environment=WANCTL_EXTERNAL_METRICS_DB=/var/lib/wanctl/metrics-spectrum.db
Environment=WANCTL_EXTERNAL_BASELINE_RTT=22.535852814520855
Environment=WANCTL_EXTERNAL_PING_SOURCE_IP=10.10.110.223
Environment=WANCTL_EXTERNAL_PING_HOSTS=1.1.1.1,9.9.9.9,208.67.222.222
Environment=CAKE_AUTORATE_BRIDGE_HEALTH_HOST=10.10.110.223
Environment=CAKE_AUTORATE_BRIDGE_HEALTH_PORT=9101
ExecStart=/usr/local/sbin/cake-autorate-spectrum-state-bridge
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
