# deploy/systemd/wanctl-nic-tuning.service
[Unit]
Description=wanctl NIC tuning (ring buffers + GRO forwarding + IRQ affinity + sysctl)
Documentation=https://github.com/kevinb361/wanctl
After=systemd-networkd-wait-online.service
Before=wanctl@spectrum.service wanctl@att.service steering.service wanctl-bridge-qos.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/wanctl-nic-tuning.sh

# Hardening (limited for root oneshot needing hardware access)
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelModules=yes
ProtectKernelTunables=no
ProtectKernelLogs=yes
ProtectClock=yes
NoNewPrivileges=no
RestrictSUIDSGID=yes
SystemCallArchitectures=native
UMask=0027

[Install]
WantedBy=multi-user.target
