# systemd/cake-continuous.timer
[Unit]
Description=CAKE Continuous Autorate - All Queues (Dual-WAN)
Documentation=https://github.com/kevin/CAKE

[Timer]
OnBootSec=30s
OnUnitActiveSec=5s
AccuracySec=1s

[Install]
WantedBy=timers.target
