# systemd/cake-att-passive.service
[Unit]
Description=CAKE Passive Queue Monitor - ATT
After=network-online.target

[Service]
Type=oneshot
User=kevin
WorkingDirectory=/home/kevin/wanctl
ExecStart=/usr/bin/python3 /home/kevin/wanctl/passive_monitor.py --config /home/kevin/wanctl/configs/att_config.yaml
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cake-att-passive

[Install]
WantedBy=multi-user.target
