# systemd/cake-spectrum-passive.service
[Unit]
Description=CAKE Passive Queue Monitor - Spectrum
After=network-online.target

[Service]
Type=oneshot
User=kevin
WorkingDirectory=/home/kevin/fusion_cake
ExecStart=/usr/bin/python3 /home/kevin/fusion_cake/passive_monitor.py --config /home/kevin/fusion_cake/configs/spectrum_config.yaml
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cake-spectrum-passive

[Install]
WantedBy=multi-user.target
