# systemd/cake-spectrum-reset.service
[Unit]
Description=Nightly CAKE Reset for Spectrum
After=network-online.target

[Service]
Type=oneshot
User=kevin
WorkingDirectory=/home/kevin/wanctl
ExecStart=/usr/bin/python3 /home/kevin/wanctl/adaptive_cake.py --config /home/kevin/wanctl/configs/spectrum_config.yaml --reset
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cake-spectrum-reset

[Install]
WantedBy=multi-user.target
