# systemd/cake-att-reset.service
[Unit]
Description=Nightly CAKE Reset for ATT
After=network-online.target

[Service]
Type=oneshot
User=kevin
WorkingDirectory=/home/kevin/fusion_cake
ExecStart=/usr/bin/python3 /home/kevin/fusion_cake/adaptive_cake.py --config /home/kevin/fusion_cake/configs/att_config.yaml --reset
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cake-att-reset

[Install]
WantedBy=multi-user.target
