# systemd/capacity-mapper.service
[Unit]
Description=CAKE Capacity Mapper - Hourly Spectrum Test
Documentation=https://github.com/kevin/CAKE

[Service]
Type=oneshot
User=kevin
WorkingDirectory=/home/kevin/wanctl
ExecStart=/home/kevin/wanctl/capacity_mapper.sh
StandardOutput=journal
StandardError=journal

# Timeout after 5 minutes (iperf + overhead)
TimeoutStartSec=300
