# systemd/steering.service
[Unit]
Description=wanctl Adaptive WAN Steering Daemon
Documentation=https://github.com/kevinb361/wanctl
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=wanctl
Group=wanctl

# FHS paths
WorkingDirectory=/opt/wanctl
ExecStart=/usr/bin/python3 -m wanctl.steering.daemon --config /etc/wanctl/steering.yaml

# Environment
Environment=PYTHONPATH=/opt
Environment=WANCTL_STATE_DIR=/var/lib/wanctl
Environment=WANCTL_LOG_DIR=/var/log/wanctl
Environment=WANCTL_RUN_DIR=/run/wanctl
EnvironmentFile=-/etc/wanctl/secrets

# Output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wanctl-steering

# Capabilities for ping (ICMP raw sockets)
AmbientCapabilities=CAP_NET_RAW

# Security hardening
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/wanctl /var/log/wanctl /run/wanctl
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
