# systemd/wanctl@.service
[Unit]
Description=wanctl Adaptive CAKE Autorate Daemon - %i
Documentation=https://github.com/kevinb361/wanctl
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wanctl
Group=wanctl

# FHS paths - %i is the WAN instance name (wan1, wan2, etc.)
WorkingDirectory=/opt/wanctl
ExecStart=/usr/bin/python3 /opt/wanctl/autorate_continuous.py --config /etc/wanctl/%i.yaml --debug

# Daemon lifecycle
Restart=on-failure
RestartSec=5s
# Watchdog: daemon pings every 2s cycle, 30s timeout gives 15x margin
WatchdogSec=30s

# Environment
Environment=PYTHONPATH=/opt
Environment=WANCTL_STATE_DIR=/var/lib/wanctl
Environment=WANCTL_LOG_DIR=/var/log/wanctl
Environment=WANCTL_RUN_DIR=/run/wanctl
EnvironmentFile=-/etc/wanctl/secrets

# Output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wanctl-%i

# Capabilities for ping (ICMP raw sockets)
AmbientCapabilities=CAP_NET_RAW

# Security hardening
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/wanctl /var/log/wanctl /run/wanctl
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
